# supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=service tomcat7 start
